# indicator.crontab
0 23 * * * source /var/www/indicator/indicator_env/bin/activate; export DB_BASE_DIR=/run/gunicorn; python /var/www/indicator/crawler.py -d $(date -I) 
