# indicator.crontab
DB_BASE_DIR=/run/gunicorn 
5 * * * * /var/www/indicator/indicator_env/bin/python /var/www/indicator/crawler.py -d $(date -I) 
