# indicator.crontab
0 22 * * * source /var/www/indicator/indicator_env/bin/activate ; python /var/www/indicator/crawler.py -d (date "+%Y-%m-%d") 
